# deploy/systemd/openslt-api.service
[Unit]
Description=OpenSLT FastAPI service
After=network.target mysql.service redis.service
[Service]
Type=simple
User=openslt
Group=openslt
WorkingDirectory=/opt/openslt
EnvironmentFile=/etc/openslt/openslt.env
ExecStart=/opt/openslt/.venv/bin/uvicorn app.main:app --app-dir backend --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
PrivateTmp=true
NoNewPrivileges=true
[Install]
WantedBy=multi-user.target

