# deploy/systemd/openslt-worker.service
[Unit]
Description=OpenSLT Celery worker
After=network.target mysql.service redis.service
[Service]
Type=simple
User=openslt
Group=openslt
WorkingDirectory=/opt/openslt
EnvironmentFile=/etc/openslt/openslt.env
Environment=PYTHONPATH=/opt/openslt/backend
ExecStart=/opt/openslt/.venv/bin/celery -A app.tasks.celery_app:celery_app worker --loglevel=INFO
Restart=on-failure
PrivateTmp=true
NoNewPrivileges=true
[Install]
WantedBy=multi-user.target
